# Charge/Email.cron
# Mailer for alpha-sb-charge
*/10 * * * * /home/alpha-sb/charge/bin/Email/doreport
